# srcs/php/src/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command = sh -c nginx -g 'daemon off;'
autostart = true
autorestart = true

[program:php7-fpm]
command = sh -c /usr/sbin/php-fpm7 --nodaemonize
autostart = true
autorestart = true